# prime/etc/nanolin.d/run.startup.cron
#!/bin/sh

#start cron
crond -l 2 -L /var/log/crontabs
