# prime/etc/nanolin.d/run.startup.cron
#!/bin/sh

[ -n "$NLDBG" ] && echo "DEBUG: $0 $@"

#start cron
crond -l 2 -L /var/log/crontabs
